# systemd_items/tempsensor_pull.timer
[Unit]
Description=Run my script every week

[Timer]
OnCalendar=weekly

[Install]
WantedBy=timers.target